# supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=docker-php-entrypoint php-fpm --nodaemonize
autostart=true
autorestart=true
priority=5

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
priority=10
stderr_logfile=/var/log/nginx.err.log
stdout_logfile=/var/log/nginx.out.log
